# .config/systemd/user/mpd.service
[Service]
Type=forking
PIDFile=/home/emallson/.config/mpd/pid
ExecStart=/usr/bin/mpd

[Install]
WantedBy=default.target