# touchegg/.config/systemd/user/touchegg@.service
[Service]
ExecStart=/usr/bin/touchegg
Environment=DISPLAY=:%i

[Install]
WantedBy=default.target