# redshift/.config/systemd/user/redshift@.service
[Service]
ExecStart=/usr/bin/redshift -l 38:-85
Environment=DISPLAY=:%i

[Install]
WantedBy=default.target